# fastapiv4.service
[Unit]
Description=FastAPI Control Server
After=network.target

[Service]
User=admin
WorkingDirectory=/home/admin/fastapi-worker

ExecStart=/home/admin/fastapi-worker/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000

Restart=always
RestartSec=3
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
